# quant_system.service
[Unit]
Description=Quant System 24/7 Server
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/quant_trading
ExecStart=/home/ec2-user/venv/bin/python -X utf8 /home/ec2-user/quant_trading/server.py
Restart=on-failure
RestartSec=10

# 환경변수 (또는 .env 파일 사용)
EnvironmentFile=/home/ec2-user/quant_trading/.env

StandardOutput=append:/home/ec2-user/quant_trading/logs/server.log
StandardError=append:/home/ec2-user/quant_trading/logs/server.log

[Install]
WantedBy=multi-user.target
